# packages/xroad-monitor-collector/redhat/SOURCES/xroad-monitor-collector.cron
#!/bin/sh
0 15 2 * * root systemctl restart xroad-monitor-collector
