# packages/xroad-monitor-collector/redhat/SOURCES/xroad-monitor-collector.cron
#!/bin/sh
HOME=/tmp
* * * * * root systemctl start xroad-monitor-collector 2>&1 > /tmp/confproxy.last;