# packages/xroad-monitor-collector/redhat/SOURCES/xroad-monitor-collector.cron
#!/bin/sh
HOME=/tmp
* * * * * xroad /usr/share/xroad/bin/xroad-monitor-collector 2>&1 > /tmp/confproxy.last;
