# packages/xroad-monitor-collector/redhat/SOURCES/xroad-monitor-collector.cron
#!/bin/sh
15 8 * * * root systemctl restart xroad-monitor-collector
