# packages/xroad-monitor-collector/redhat/SOURCES/xroad-monitor-collector.cron
#!/bin/sh
* * * * * root systemctl restart xroad-monitor-collector
